# logger.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/home/stu/tacoma_data_logger/bin/python3 /home/stu/tacoma_data_logger/battery_logger.py
Restart=on-failure
RestartSec=60
User=stu
WorkingDirectory=/home/stu/tacoma_data_logger/
Environment="PYTHONBUFFERED=1"

[Install]
WantedBy=multi-user.target
