# logger.service
[Unit]
Description=Python script for logging data from epever controller
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncdb
User=stu
WorkingDirectory=/usr/local/bin/

[Install]
WantedBy=multi-user.target
