# syncdb.timer
[Unit]
Description=timer for syncdb

[Timer]
OnUnitActiveSec=30s
OnBootSec=10s

[Install]
WantedBy=timers.target
